# supervisord.conf
; supervisord.conf
;
; Runs gunicorn + a Celery worker + Celery beat inside the single Render web
; container. Logs go to stdout/stderr so they show up in Render's log stream.
;
; Free-tier memory notes (512 MB shared):
;   * gunicorn: 1 gthread worker (1 process, 4 threads) — one app copy in RAM.
;   * celery worker: solo pool, concurrency 1 — one app copy in RAM.
;   * celery beat: tiny scheduler process.
; Watch Render's memory graph after deploy; if tight, the first lever is to
; move bot work to the cron-triggered inline path and drop the worker.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
; supervisord itself runs as PID 1 (via entrypoint exec) and reaps children.

[program:gunicorn]
command=daphne -b 0.0.0.0 -p %(ENV_PORT)s config.asgi:application
directory=/app
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:celery-worker]
command=celery -A config worker --loglevel=info --pool=solo --concurrency=1 --without-gossip --without-mingle --without-heartbeat
directory=/app
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:celery-beat]
command=celery -A config beat --loglevel=info --schedule=/tmp/celerybeat-schedule
directory=/app
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
